# listenfifo.socket
[Socket]
ListenFIFO=/tmp/listenfifo

[Install]
WantedBy=sockets.target

