# listenfifo.service
[Unit]
Description=Listen special

[Service]
ExecStart=/usr/bin/env python3 /home/pi/listenfifo/listenfifo.py

